# VPS/install/monitoring/etc/systemd/vps-local-metrics.timer
[Unit]
Description=Publication périodique des métriques locales du VPS

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m
AccuracySec=1m
Persistent=false
Unit=vps-local-metrics.service

[Install]
WantedBy=timers.target
