# VPS/install/system/systemd/vps-nightly-maintenance.service
[Unit]
Description=Sauvegarde et maintenance nocturne du VPS
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/vps-maintenance
ExecStart=/usr/local/sbin/vps-nightly-maintenance
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUWeight=20
RuntimeMaxSec=3h30m
